# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fcgiwrap]
command=/usr/sbin/fcgiwrap -s unix:/var/run/fcgiwrap.socket
socket=unix:///var/run/fcgiwrap.socket
socket_owner=www-data:www-data
socket_mode=0770
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/fcgiwrap.log
stderr_logfile=/var/log/supervisor/fcgiwrap.log

[program:smokeping]
command=/usr/sbin/smokeping --config=/etc/smokeping/config --nodaemon
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/smokeping.log
stderr_logfile=/var/log/supervisor/smokeping.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log